# deploy/hostinger/elore-paris.service
[Unit]
Description=Elore Paris Next.js storefront
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=elore
Group=elore
WorkingDirectory=/srv/elore-paris/current
EnvironmentFile=/etc/elore-paris/elore-paris.env
Environment=NODE_ENV=production
Environment=HOSTNAME=127.0.0.1
Environment=PORT=3056
Environment=COZMATEKS_PROJECT_ROOT=/srv/elore-paris/current
ExecStart=/usr/bin/node /srv/elore-paris/current/.next/standalone/server.js
Restart=on-failure
RestartSec=5
TimeoutStopSec=30
KillSignal=SIGTERM
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
CacheDirectory=elore-paris
CacheDirectoryMode=0750
ReadWritePaths=/var/lib/elore-paris
UMask=0027

[Install]
WantedBy=multi-user.target
